# systemd/xueqiu-scraper.timer
[Unit]
Description=Run Xueqiu Article Scraper daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
